# rsrc_agent/supervisord.conf
[supervisord]
nodaemon=true
[program:mms]
command=acronis_mms
[program:crond]
command=cron -f