# rsrc_server/prepare_config.service
[Service]
Type=oneshot
ExecStart=/usr/local/libexec/prepare_config.sh
RemainAfterExit=yes
Before=basic.target

[Install]
RequiredBy=basic.target